# planner-web.service
[Unit]
Description=Educational Schedule Planner Web Service
After=network.target

[Service]
User=user
Group=user
WorkingDirectory=/home/user/planner-solving
Environment="PATH=/home/user/.pyenv/shims:/home/user/.pyenv/bin:$PATH"
Environment="PATH=/home/user/.pyenv/shims:/home/user/.pyenv/bin:$PATH"
Environment="PYENV_ROOT=/home/user/.pyenv"
Environment="PYTHONPATH=/home/user/planner-solving"
ExecStart=/home/user/.pyenv/shims/python -m uvicorn web.backend.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target